# package/etc/cron.d/gitaptly
0 * * * * root /bin/bash /usr/bin/gitaptly_update.sh

