# package/etc/cron.d/gitaptly
0 * * * * root /usr/share/gitaptly/update

