# package/etc/cron.d/gitaptly
0 * * * * root /bin/sh -c 'PATH=/usr/share/gitaptly:"$PATH"; exec update'
