# package/DEBIAN/postinst
#!/bin/bash -e

if [ "$1" = 'abort-configure' ] || [ "$1" = 'abort-upgrade' ] || [ "$1" = 'abort-remove' ] || [ "$1" = 'abort-deconfigure' ]; then
  exit 0
elif [ "$1" != 'configure' ] && [ "$1" != 'reconfigure' ]; then
  exit 1
fi

rm -rf /opt/gitaptly/* || true
rm -rf /var/lib/gitaptly/cgi-bin && sed "s/cgi-bin/pool/" /var/lib/gitaptly/dists/stable/main/binary-all/Packages || true # legacy


. /usr/share/debconf/confmodule
db_get gitaptly/PORT
echo "PORT=$RET" >> /opt/gitaptly/env
db_get gitaptly/MODE
echo "MODE=$RET" >> /opt/gitaptly/env
db_get gitaptly/GITHUB_API_TOKEN
echo "GITHUB_API_TOKEN=\"$RET\"" >> /opt/gitaptly/env
db_get gitaptly/OTLP_TRACES_ENDPOINT
echo "OTLP_TRACES_ENDPOINT=\"$RET\"" >> /opt/gitaptly/env
db_get gitaptly/OTLP_TRACES_HEADER
echo "OTLP_TRACES_HEADER=\"$RET\"" >> /opt/gitaptly/env
db_get gitaptly/OTLP_METRICS_ENDPOINT
echo "OTLP_METRICS_ENDPOINT=\"$RET\"" >> /opt/gitaptly/env
db_get gitaptly/OTLP_METRICS_HEADER
echo "OTLP_METRICS_HEADER=\"$RET\"" >> /opt/gitaptly/env
db_get gitaptly/OTLP_LOGS_ENDPOINT
echo "OTLP_LOGS_ENDPOINT=\"$RET\"" >> /opt/gitaptly/env
db_get gitaptly/OTLP_LOGS_HEADER
echo "OTLP_LOGS_HEADER=\"$RET\"" >> /opt/gitaptly/env
db_stop

source /opt/gitaptly/env

touch /var/lib/gitaptly/pool
if [ "$MODE" = "cache" ]; then
  if ! [ -d /var/lib/gitaptly/pool ]; then rm -rf /var/lib/gitaptly/pool && mkdir -p /var/lib/gitaptly/pool; fi
elif [ "$MODE" = "proxy" ]; then
  rm -rf /var/lib/gitaptly/pool || true
  ln --symbolic /usr/share/gitaptly/serve /var/lib/gitaptly/pool
fi
bash /usr/share/gitaptly/update

systemctl daemon-reload
systemctl enable gitaptly
systemctl start gitaptly
