# package/DEBIAN/prerm
#!/bin/bash
systemctl stop gitaptly
systemctl disable gitaptly
exit 0
